# ext/systemd/jerakia.service
[Unit]
Description=Jerakia Server

[Service]
ExecStart=/opt/jerakia/bin/jerakia server
Type=simple


[Install]
WantedBy=multi-user.target
